# .config/systemd/user/poweralertd.service
[Unit]
Description=poweralertd

[Service]
Type=simple
ExecStart=/usr/bin/poweralertd

[Install]
WantedBy=default.target
