# .config/systemd/user/wlsunset.service
[Unit]
Description=wlsunset

[Service]
Type=simple
ExecStart=/home/rai/bin/wlsunset

[Install]
WantedBy=default.target
